# os/systemd/valkey.service.d/homeharbor.conf
[Service]
RuntimeDirectory=valkey
RuntimeDirectoryMode=0750
UMask=0007
ExecStart=
ExecStart=/usr/bin/valkey-server /etc/valkey/valkey.conf --supervised systemd --daemonize no --unixsocket /run/valkey/homeharbor.sock --unixsocketperm 770 --port 0
