# os/systemd/audit-rules.service.d/homeharbor.conf
[Service]
ExecStart=
ExecStart=/usr/bin/auditctl -R /etc/audit/homeharbor.rules
