# os/systemd/homeharbor-container-apply.path
[Unit]
Description=Watch HomeHarbor container apply requests

[Path]
PathChanged=/run/homeharbor/container-apply.request
Unit=homeharbor-container-apply.service

[Install]
WantedBy=multi-user.target
