# os/systemd/homeharbor-container-apply.service
[Unit]
Description=Apply HomeHarbor rootless container Quadlet configuration
Requires=homeharbor-api.service
After=network-online.target homeharbor-api.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/lib/homeharbor/agent/HomeHarbor.Agent apply-containers
