# os/systemd/homeharbor-boot-attempt.service
[Unit]
Description=Record HomeHarbor boot attempt and enter recovery after boot loops
DefaultDependencies=no
After=local-fs.target
Before=multi-user.target
ConditionPathExists=/usr/lib/homeharbor/agent/HomeHarbor.Agent

[Service]
Type=oneshot
ExecStart=/usr/lib/homeharbor/agent/HomeHarbor.Agent boot-attempt

[Install]
WantedBy=multi-user.target
