# os/systemd/homeharbor-storage-health.timer
[Unit]
Description=Run HomeHarbor storage health check hourly

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target

