# os/systemd/homeharbor-caddy-render.service
[Unit]
Description=Render HomeHarbor Caddy reverse proxy config
Requires=homeharbor-api.service
After=network-online.target homeharbor-api.service
Wants=network-online.target
StartLimitIntervalSec=5min
StartLimitBurst=6

[Service]
Type=oneshot
ExecStart=/usr/lib/homeharbor/agent/HomeHarbor.Agent render-caddyfile
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
