# os/systemd/homeharbor-recovery-action.path
[Unit]
Description=Watch for physical HomeHarbor recovery console actions

[Path]
PathChanged=/run/homeharbor-recovery/action.request
Unit=homeharbor-recovery-action.service

[Install]
WantedBy=multi-user.target
