# os/systemd/homeharbor-boot-success.service
[Unit]
Description=Mark HomeHarbor boot successful after control plane health is ready
Requires=homeharbor-selinux-ready.target homeharbor-api.service
After=homeharbor-selinux-ready.target multi-user.target homeharbor-api.service network-online.target
Wants=network-online.target
ConditionPathExists=/usr/lib/homeharbor/agent/HomeHarbor.Agent

[Service]
Type=oneshot
ExecStart=/usr/lib/homeharbor/agent/HomeHarbor.Agent boot-success

[Install]
WantedBy=multi-user.target
